# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh
umask 0002

exec php-fpm --nodaemonize;
